# hbase/services/rest.service.conf
[program:rest]
command=hbase rest start
stdout_events_enabled = true
stderr_events_enabled = true
