# hbase/services/thrift.service.conf
[program:thrift]
command=hbase thrift start -t 5000 -p 9090
stdout_events_enabled = true
stderr_events_enabled = true
