# hbase/services/zookeeper.service.conf
[program:zookeeper]
command=hbase zookeeper start
stdout_events_enabled = true
stderr_events_enabled = true
